# ansible/roles/common/templates/devmon.service.j2
[Unit]
Description=Automounting usb drives.
After=network.target

[Service]
Type=simple
User=pi
Restart=on-abort
ExecStart=/usr/bin/devmon

[Install]
WantedBy=multi-user.target