# ansible/roles/common/templates/hs-iptables.service.j2
[Unit]
Description=Activate IPtables for Hotspot
After=network-pre.target
Before=network-online.target

[Service]
Type=simple
ExecStart=/etc/iptables-rules.v4

[Install]
WantedBy=multi-user.target
