# ansible/roles/common/templates/scratch-gui.service.j2
[Unit]
Description=Scratch GUI
After=network.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/npm start --prefix /opt/scratch-gui
Restart=on-failure

[Install]
WantedBy=multi-user.target
